# crystalbay.service
[Unit]
Description=Crystal Bay Travel Application
After=network.target postgresql.service

[Service]
Type=exec
User=crystalbay
Group=crystalbay
WorkingDirectory=/home/crystalbay/crystal-bay-travel
Environment=PATH=/home/crystalbay/crystal-bay-travel/venv/bin
EnvironmentFile=/home/crystalbay/crystal-bay-travel/.env
ExecStart=/home/crystalbay/crystal-bay-travel/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 3 --timeout 120 --keep-alive 2 --max-requests 1000 --preload main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target